# build/faas-build/supervisord.conf
[supervisord]
nodaemon=true

[program:faas-api]
command=/root/.local/bin/faas-api
autorestart=true
startretries=50
priority=0
